# configs/systemd/modem-watchdog.timer.template
# p5g watchdog timer.
# Installed to: /etc/systemd/system/p5r-watchdog.timer
#
# Fires the watchdog service every HEALTHCHECK_INTERVAL seconds.

[Unit]
Description=p5g watchdog timer

[Timer]
OnBootSec=2min
OnUnitActiveSec=60s
Unit=p5r-watchdog.service

[Install]
WantedBy=timers.target
